# install/files/hassio-apparmor.conf
description=Hass.io AppArmor
start on syno.network.ready
stop on runlevel [06]
exec /usr/sbin/hassio-apparmor
