# install/files/hassio-supervisor.conf
description=Hass.io supervisor
start on syno.network.ready
stop on runlevel [06]
exec /usr/sbin/hassio-supervisor
